# ansible/production/roles/backup/templates/backup-daily.timer.j2
[Unit]
Description=Perform daily backup

[Timer]
OnCalendar=daily
RandomizedDelaySec=3600
Unit=backup.service

[Install]
WantedBy=timers.target
