# ansible/production/roles/backup/templates/backup.service.j2
[Unit]
Description=Perform backup using /root/backup.sh


[Service]
Type=simple
Nice=19
IOSchedulingClass=2
IOSchedulingPriority=7
ExecStart=/root/backup.sh
